# deployment/system-service/dist-contents/alis.service
[Unit]
Description=ALIS Service
After=syslog.target network.target

[Service]
Type=simple
PIDFile=/u01/LIMS_APi/alis/RUNNING_PID
#EnvironmentFile=/u01/LIMS_APi/alis/alis-env.conf
User=alis
WorkingDirectory=/u01/LIMS_APi/alis/
ExecStart=/bin/java -Xms2g -Xmx2g -jar alis.jar

Restart=on-abort

[Install]
WantedBy=multi-user.target

